# systemd/.config/systemd/user/change_wallpaper.service
[Unit]
Description=Change Wallpaper Service

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/change_wallpaper.sh

